# files/etc/init.d/hacks
#!/sbin/openrc-run
# unmounts /run and points in to /var/run

description="Fix run"

start() {
    ebegin "Running hacks to make system fit for flash memory systems"
    if [ ! -L "/run/run" ]; then
        eindent
        ebegin "mounting /run on /var/run"
        ln -s /run /run
        eend 0
        eoutdent
    fi
    eend 0
}
